# salt/mezzanine/production/files/production_supervisor.conf
{% set settings = pillar.get('settings', {}) -%}
{% set env = settings.get('env', {}) -%}
{% set proj_name = env.get('proj_name', '') -%}
{% set proj_path = env.get('proj_path', '') -%}
{% set venv_path = env.get('venv_path', '') -%}
{% set user = env.get('user', '') -%}

[program:uwsgi_{{ proj_name }}]
directory={{ proj_path }}
command={{ venv_path }}/bin/uwsgi /home/django/uwsgi.ini
user={{ user }}
stopsignal=QUIT
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=20
